# nix/docker/all_in_one/init/s6-services/services/acme-selfsigned/up
#!/bin/execlineb -P

/etc/s6-services/acme-selfsigned/shell_up.sh
