# nix/docker/all_in_one/init/s6-services/services/nocodb-main-log/run
#!/bin/execlineb -P

envfile /etc/s6-confs/nocodb.main.conf

importas -sCuD "n3 s2000000 T" DIRECTIVES DIRECTIVES

ifelse { test -w /var/log } {
	foreground { install -d -o s6log -g s6log /var/log/nocodb-main }
	s6-setuidgid s6log exec -c s6-log -d3 -b -- ${DIRECTIVES} /var/log/nocodb-main
}

foreground { install -d -o s6log -g s6log /run/log/nocodb-main }

s6-setuidgid s6log exec -c s6-log -d3 -b -- ${DIRECTIVES} /run/log/nocodb-main
