# nix/docker/all_in_one/init/s6-services/services/postgresql-db-create/up
#!/bin/execlineb -P

importas -sCuD "n3 s2000000 T" DIRECTIVES DIRECTIVES

foreground { install -d -o s6log -g s6log /var/log/postgresql-db-create }
fdmove -c 2 1

pipeline -dw {
	s6-setuidgid s6log exec -c s6-log -b -- ${DIRECTIVES} /var/log/postgresql-db-create
}

if -n {
	pipeline {
		psql -U postgres -tAc SELECT\ 1\ FROM\ pg_database\ WHERE\ datname\ =\ \'nocodb\'
	} grep -q 1
} psql -U postgres -tAc CREATE\ DATABASE\ \"nocodb\"
