# nix/docker/all_in_one/init/s6-services/services/nocodb-worker-srv/run
#!/bin/execlineb -P

# setup by ../../../env-processor
envfile /etc/s6-confs/nocodb.worker.conf
envfile /run/nocodb.common.dynamic.env
envfile /run/nocodb.worker.dynamic.env

foreground { mkdir -p /var/lib/nocodb }
foreground { chown nocodb:nocodb /tmp /var/lib/nocodb }
foreground {
	if {
		test -d /var/lib/nocodb/nc
	} chown -R nocodb:nocodb /var/lib/nocodb/nc
}
foreground {
	if {
		test -f /var/lib/nocodb/noco.db
	} chown nocodb:nocodb /var/lib/nocodb/noco.db
}

fdmove -c 2 1
execline-cd  /var/lib/nocodb s6-setuidgid nocodb nocodb
